# systemd/services/set_display_temp_normal.service
[Unit]
Description=Set gamma to 6500K

[Service]
Type=oneshot
ExecStart=%h/.dotfiles/scripts/systemd/set_temp_low.sh

